# source/debian/postinst
#!/bin/sh
# postinst script for desktop-theme-kit
#
# see: dh_installdeb(1)

set -e


THEME=kit
PACKAGE=desktop-theme-package


# Count the previously installed desktop-theme-packages in /usr/share
# and add their number to the priority. So the last installed package
# always gets the highest priority.
base_priority=80
desktop_theme_packages=$(find /usr/share/desktop-theme-* -prune | wc -l)
priority=$((base_priority+$desktop_theme_packages))

install_alternatives() {
    update-alternatives --install \
        /usr/share/images/desktop-base/desktop-background \
        desktop-background \
        /usr/share/desktop-theme-$THEME/backgrounds/$THEME-wallpaper.png $priority

	# desktop-background.xml means actually the Gnome lockscreen only.
    update-alternatives --install \
        /usr/share/images/desktop-base/desktop-background.xml \
        desktop-background.xml \
        /usr/share/desktop-theme-$THEME/backgrounds/$THEME-gnome-lockscreen.xml $priority

    update-alternatives --install \
        /usr/share/images/desktop-base/desktop-grub.png \
        desktop-grub \
        /usr/share/desktop-theme-$THEME/grub/$THEME-grub.png $priority

    update-alternatives --install \
	/usr/share/images/desktop-base/desktop-login \
	desktop-login \
	/usr/share/desktop-theme-$THEME/backgrounds/$THEME-login.png $priority
}

set_alternatives_auto() {
    update-alternatives --auto desktop-background
    update-alternatives --auto desktop-background.xml
    update-alternatives --auto desktop-grub
    update-alternatives --auto desktop-login
}

install_logfile() {
    # To keep track of the installed plymouth themes, we use a logfile in /var.
    test -f /var/log/desktop-theme-packages.log || touch /var/log/desktop-theme-packages.log
    echo "$THEME" >> /var/log/desktop-theme-packages.log
}

update_plymouth_config() {
    if which plymouth-set-default-theme > /dev/null ; then
        plymouth-set-default-theme $THEME || true
    fi
}

update_grub_config() {
    if which update-grub2 > /dev/null ; then
        sync
        update-grub2 || true
    fi

    if [ -x /usr/sbin/update-initramfs ]; then
        update-initramfs -u
    fi
}

update_gdm3_config() {
    # Reset the gdm3 wallpaper cache
    rm -f /var/lib/gdm3/.cache/wallpaper/*
    # Reload the gdm3 configuration
    if [ -x /usr/sbin/gdm3 ]; then
        invoke-rc.d gdm3 reload || true
    fi
}

update_kdm_config() {
    # The KDM theme is configured when KDM is started. To understand what's going on,
    # take a look at the function "setup_config" in "/etc/init.d/kdm".
    # So if KDM is already running, we need to reload it to make our theme visible.
    if [ -f /run/kdm.pid ]; then
        invoke-rc.d kdm reload || true
    fi
}

case "$1" in
    configure|abort-upgrade)
        install_alternatives
        set_alternatives_auto
        install_logfile
        update_grub_config
        update_gdm3_config
        update_kdm_config
        update_plymouth_config
    ;;

    abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
